# packaging/debian/postinst
#!/bin/bash
set -e

# Reload systemd user units
su "$SUDO_USER" -c "systemctl --user daemon-reload"

# Enable service for the current user (optional but recommended)
su "$SUDO_USER" -c "systemctl --user enable sonik.service" || true

# Start daemon immediately
su "$SUDO_USER" -c "systemctl --user start sonik.service" || true

exit 0
